# deploy/casinonet-dhcp.service
[Unit]
Description=CabiNet DHCP server (slot network, pinned to eth0)
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/home/aj/CasinoNet/G2S
ExecStart=/usr/bin/python3 -u start-dhcp-enhanced.py --interface eth0
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
