# Companion/casinonet-companion.service
[Unit]
Description=CabiNet Companion (PN532 RFID taps -> hub, rides beside casinonet-sas)
After=network-online.target
Wants=network-online.target
# GR-20 storm cap, same rationale as casinonet-sas.service. companion_host.py
# tolerates a missing/mute PN532 INTERNALLY (readerOk=false, lazy retry — the
# unit stays active), so a board that isn't wired yet never trips this limit.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=aj
Group=aj
# /dev/i2c-11 (the software i2c-gpio bus on GPIO23/24) is root:i2c; aj is a
# member, SupplementaryGroups makes the unit robust if login groups change.
SupplementaryGroups=i2c
WorkingDirectory=/home/aj/CasinoNet/Companion
# Stdlib-only on purpose — system python3 is enough on the Zero (no venv).
# ZERO-CONFIG (v11): FLAGLESS by design so ONE identical image boots on every
# card and self-configures — the daemon IDs itself by the Pi's hardware serial
# (companion-<serial-tail>) and finds the hub at the DHCP default gateway (the
# host is the gateway + DHCP server on the wired slot segment; satellites are
# wired-only). The machine binding is then ASSIGNED FROM THE HUB UI (Players ▸ Readers ▸
# Assign) — the hub routes taps by companionId, so no --g2s-egm here and a
# re-assign takes effect live with no restart. companion_setup.sh can still
# inject explicit flags (--hub / --companion-id / --g2s-egm / --sas-*) for a
# co-located or manual override; a bare flag-free start is the collector path.
ExecStart=/usr/bin/python3 -u /home/aj/CasinoNet/Companion/companion_host.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
